# preciso.service
[Unit]
Description=Preciso API
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/opt/preciso
Environment=APP_ENV=prod
Environment=PUBLIC_DOMAIN=preciso-data.com
Environment=FINROBOT_ENABLED=1
Environment=FINROBOT_OAI_CONFIG=/opt/preciso/OAI_CONFIG_LIST
Environment=FINROBOT_API_KEYS=/opt/preciso/config_api_keys
Environment=FINROBOT_MODEL=gpt-4-0125-preview
ExecStart=/opt/preciso/.venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8004
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
